# custom_steps/etcd/templates/etcd.service.j2
[Unit]
Description=etcd key-value store
Documentation=https://github.com/etcd-io/etcd
After=network.target

[Service]
User=root
Type=notify
EnvironmentFile=-/etc/default/etcd
ExecStart=/usr/local/bin/etcd \
  --name {{ etcd_name | default(inventory_hostname) }} \
  --data-dir {{ etcd_data_dir }} \
  --listen-client-urls {% if etcd_enable_tls %}https{% else %}http{% endif %}://0.0.0.0:{{ etcd_client_port }} \
  --advertise-client-urls {% if etcd_enable_tls %}https{% else %}http{% endif %}://{{ ansible_host }}:{{ etcd_client_port }} \
  --listen-peer-urls {% if etcd_enable_tls %}https{% else %}http{% endif %}://0.0.0.0:{{ etcd_peer_port }} \
  --initial-advertise-peer-urls {% if etcd_enable_tls %}https{% else %}http{% endif %}://{{ ansible_host }}:{{ etcd_peer_port }} \
  --initial-cluster-token {{ etcd_cluster_token }} \
  --initial-cluster "{% for host in groups['etcd_cluster'] %}{{ hostvars[host].etcd_name | default(host) }}={% if etcd_enable_tls %}https{% else %}http{% endif %}://{{ hostvars[host].ansible_host }}:{{ etcd_peer_port }}{{ ',' if not loop.last }}{% endfor %}" \
  --initial-cluster-state new \
  --enable-v2 {{ etcd_enable_v2 | lower }} \
  {% if etcd_enable_tls %} \
  --client-cert-auth \
  --trusted-ca-file {{ etcd_tls_dir }}/{{ etcd_ca_cert }} \
  --cert-file {{ etcd_tls_dir }}/{{ etcd_server_cert }} \
  --key-file {{ etcd_tls_dir }}/{{ etcd_server_key }} \
  --peer-client-cert-auth \
  --peer-trusted-ca-file {{ etcd_tls_dir }}/{{ etcd_ca_cert }} \
  --peer-cert-file {{ etcd_tls_dir }}/{{ etcd_server_cert }} \
  --peer-key-file {{ etcd_tls_dir }}/{{ etcd_server_key }} \
  {% endif %}
Restart=always
RestartSec=5
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target